# kiosk-web.service
[Unit]
Description=Kiosk Web Manager
After=network.target

[Service]
Type=simple
User=annkiosk
WorkingDirectory=/home/annkiosk/announcements_kiosk
ExecStart=/home/annkiosk/announcements_kiosk/venv/bin/python3 /home/annkiosk/announcements_kiosk/web_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
